# com.ap.blocker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ap.blocker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/blocker</string>
    </array>

    <!-- Start immediately when agent is loaded (at login) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart automatically if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Minimum seconds between restarts to avoid spin-loop on crash -->
    <key>ThrottleInterval</key>
    <integer>10</integer>
</dict>
</plist>
